# usbutler/supervisord-web.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord-web.log
pidfile=/var/run/supervisord-web.pid
user=root

[program:pcscd]
command=/usr/sbin/pcscd --foreground --disable-polkit
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/pcscd-web.err.log
stdout_logfile=/var/log/supervisor/pcscd-web.out.log
priority=100

[program:web]
command=/app/.venv/bin/python -m app.web.server
directory=/app
user=appuser
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/web.err.log
stdout_logfile=/var/log/supervisor/web.out.log
environment=PATH="/app/.venv/bin:%(ENV_PATH)s"
priority=200
